# .docker/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command = /etc/init.d/apache2 start

[program:mysql]
command = /usr/bin/mysqld_safe